# packaging/bm-webmail/ROOT/etc/cron.d/bm-webmail
0 */6 * * * root [ -d /usr/share/bm-webmail/temp/ ] && find /usr/share/bm-webmail/temp -type f -cmin +360 -exec rm -f {} \; 
